# templates/hadoop_exporter.service
[Unit]
Description=Hadoop Prometheus Exporter (offline bundle, kein Docker)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Pfad anpassen: dorthin, wohin das Bundle auf dem Zielserver entpackt wurde
WorkingDirectory=%h/hadoop_exporter
ExecStart=%h/hadoop_exporter/start.sh
Restart=on-failure
RestartSec=5s

NoNewPrivileges=true
ProtectSystem=full

[Install]
WantedBy=default.target
